# cron/crontab
* * * * * python3.6 /opt/app/src/main.py >> /var/log/cron.log 2>&1
